# dump-database.timer
#configuration du Timer
[Unit]
Description=10min
[Timer]
OnCalendar= 0:10
[Install]
WantedBy=timers.target